# deployments/service/ignis-webhook.service
[Unit]
Description=Ignis GitHub Webhook Listener
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ignis/deployments/webhook
ExecStart=/usr/local/bin/bun server.ts >> /opt/ignis/logs/webhook/service.log 2>> /opt/ignis/logs/webhook/service-error.log
EnvironmentFile=/opt/ignis/.env
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
