# deployments/service/ignis-webhook.service
[Unit]
Description=Ignis GitHub Webhook Listener
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ignis/deployments/webhook
ExecStart=/usr/local/bin/bun run server.ts
Environment=NODE_ENV=production
EnvironmentFile=/opt/ignis/.env
Restart=on-failure
RestartSec=3
StandardOutput=append:/opt/ignis/logs/webhook/service.log
StandardError=append:/opt/ignis/logs/webhook/service-error.log

[Install]
WantedBy=multi-user.target
