# deployments/service/ignis-webhook.service
[Unit]
Description=Ignis GitHub Webhook Listener
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=ignis
Group=ignis
WorkingDirectory=/opt/ignis/deployments/webhook
ExecStart=/usr/local/bin/bun run server.ts
StandardOutput=journal
StandardError=journal
EnvironmentFile=/opt/ignis/.env
Restart=on-failure
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30

PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target