# deployments/service/ignis-webhook.service
[Unit]
Description=Ignis GitHub Webhook Listener
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=ignis
Group=ignis
WorkingDirectory=/opt/ignis/deployments/webhook
ExecStart=/bin/bash -c '/usr/local/bin/bun run server.ts'
StandardOutput=file:/opt/ignis/logs/webhook/webhook-$(date +\%Y\%m\%d-\%H\%M\%S).log
StandardError=file:/opt/ignis/logs/webhook/webhook-error-$(date +\%Y\%m\%d-\%H\%M\%S).log
EnvironmentFile=/opt/ignis/.env
Restart=on-failure
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30

PrivateTmp=true
NoNewPrivileges=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target