# deployments/service/ignis-startup.service
[Unit]
Description=Ignis System Startup
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/ignis/deployments/scripts/startup.sh
WorkingDirectory=/opt/ignis
EnvironmentFile=/opt/ignis/.env
StandardOutput=append:/opt/ignis/logs/startup/service.log
StandardError=append:/opt/ignis/logs/startup/service-error.log
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
